# ansible/templates/pi-analytics-kiosk.service.j2
[Unit]
Description=Pi Analytics Dashboard Kiosk Mode
After=multi-user.target graphical.target pi-analytics-backend.service getty@tty1.service
Wants=graphical.target network-online.target

[Service]
Type=simple
User={{ ansible_user }}
Group={{ ansible_user }}
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/{{ ansible_user }}/.Xauthority"
Environment="HOME=/home/{{ ansible_user }}"
# Wait for backend and display to be ready
ExecStartPre=/bin/sleep 10
ExecStart={{ app_dir }}/start-kiosk.sh
Restart=always
RestartSec=10
StandardOutput=append:/var/log/pi-analytics/kiosk.log
StandardError=append:/var/log/pi-analytics/kiosk-error.log

[Install]
WantedBy=default.target