# ansible/templates/pi-analytics-display.service.j2
[Unit]
Description=Pi Analytics Dashboard Display Service
After=multi-user.target pi-analytics-backend.service hyperpixel2r-init.service network.target
Wants=pi-analytics-backend.service
Requires=pi-analytics-backend.service

[Service]
Type=simple
User=root
Group=root

# X server environment
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/{{ app_user }}/.Xauthority"
Environment="HOME=/home/{{ app_user }}"

# Start X server with window manager and browser
ExecStart=/usr/bin/xinit {{ app_dir }}/scripts/start-display.sh -- :0 -nocursor vt2

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Output to journal
StandardOutput=journal
StandardError=journal

# Kill any processes when stopping
KillMode=control-group

[Install]
WantedBy=multi-user.target